# config/systemd/user/syncthing.service
[Unit]
Description=Syncthing - File sync client

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
